# systemd/learnify-backup.service
[Unit]
Description=Learnify automatic verified backup to Google Drive
After=network-online.target mariadb.service apache2.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/learnify-backup.sh
