# systemd/learnify-backup.service
[Unit]
Description=Learnify backup runner
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin"
ExecStart=/usr/local/sbin/backup-learnify.sh
ExecStartPre=/usr/bin/flock -n /var/lock/learnify-backup.lock -c /bin/true
