# systemd/learnify-backup.timer
[Unit]
Description=Daily Learnify backup at 03:07 UTC
[Timer]
OnCalendar=*-*-* 03:07:00
Persistent=true
RandomizedDelaySec=90
[Install]
WantedBy=timers.target
