# etc/cron.d/backup-learnify
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

20 3 * * * root /usr/local/sbin/backup-learnify.sh >/var/log/backup-learnify.log 2>&1
