# golden/snap-20251203-135659UTC/etc/systemd/system/learnify-workflows.service.d/10-glass.conf
[Service]
Environment=GLASS_PROVIDER=fs
Environment=GLASS_FS_ROOT=/opt/learnify/worktree/app/modules/workflow_lite/work
